# src/supervisor/conf.d/queue_api_docs.conf
[program:queue_api_docs]
command=/var/www/teorema/env/bin/python3 asyncapi_server.py
directory=/var/www/teorema/queue_api/resources/asyncapi
user=www-data
autorestart=true
redirect_stderr=true

